# examples/deploy/supervisor.conf
[program:huey]
directory=/srv/my_app
; Run the consumer directly -- no shell-script wrappers -- so that signals
; are delivered to the Python process. See github issue #88.
command=/srv/my_app/venv/bin/huey_consumer my_app.huey -w 4 -k thread
user=appuser
autostart=true
autorestart=true

; Huey shuts down gracefully on SIGINT, allowing workers to finish their
; current task. Supervisor's default stopsignal is TERM, which huey treats
; as "stop immediately, interrupting any running tasks" -- so use INT.
stopsignal=INT
; How long to wait for in-flight tasks to finish before escalating to
; SIGKILL. Set this to comfortably exceed your longest-running task.
stopwaitsecs=60

stdout_logfile=/var/log/huey/stdout.log
stderr_logfile=/var/log/huey/stderr.log
environment=PYTHONUNBUFFERED="1"
